# supervisor/dev.conf
[supervisord]
nodaemon=true

[program:webserver]
stopasgroup=true
directory = /opt/personnel-testing
command = python3 ./manage.py runserver 0.0.0.0:8080